# crontab
# Contests Spider
0 * * * * /bin/sh /work/spider.sh
